# openclaw-userlook/deploy/systemd/openclaw-userlook-backend.service.example
[Unit]
Description=openclaw-userlook FastAPI backend
After=network.target mysql.service
Wants=network.target

[Service]
Type=simple
User=openclaw
Group=openclaw
WorkingDirectory=/opt/openclaw-userlook/backend
EnvironmentFile=/opt/openclaw-userlook/backend/.env
ExecStart=/opt/openclaw-userlook/backend/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 10009
Restart=always
RestartSec=5
StandardOutput=append:/opt/openclaw-userlook/backend/logs/backend.log
StandardError=append:/opt/openclaw-userlook/backend/logs/backend.err.log

[Install]
WantedBy=multi-user.target
